# csiem/windows/file/file_event/file_event_win_susp_task_write.json
{
  "Description": "Detects the creation of tasks from processes executed from suspicious locations\n\nAuthor: Florian Roth (Nextron Systems)\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "80e1f67a-4596-4351-98f5-a9c3efabac95",
  "InsertDate": "2025-01-25T21:21:03Z",
  "LastUpdateDate": "2025-01-25T21:21:03Z",
  "Level": "high",
  "Name": "Suspicious Scheduled Task Write to System32 Tasks",
  "Query": "sourcetype='windows-sysmon' eql select * from _source_ where filepath like '%\\Windows\\System32\\Tasks%' and (process_path like '%\\AppData\\%' or process_path like '%C:\\PerfLogs%' or process_path like '%\\Windows\\System32\\config\\systemprofile%')",
  "Tags": [
    "attack.persistence",
    "attack.execution",
    "attack.t1053"
  ],
  "Version": "0.1"
}