# csiem/macos/process_creation/proc_creation_macos_persistence_via_plistbuddy.json
{
  "Description": "Detects potential persistence activity using LaunchAgents or LaunchDaemons via the PlistBuddy utility\n\nAuthor: Sohan G (D4rkCiph3r)\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "65d506d3-fcfe-4071-b4b2-bcefe721bbbb",
  "InsertDate": "2025-01-25T21:19:38Z",
  "LastUpdateDate": "2025-01-25T21:19:38Z",
  "Level": "high",
  "Name": "Potential Persistence Via PlistBuddy",
  "Query": "sourcetype like 'macos-%' eql select * from _source_ where process_path like '%/PlistBuddy' and (command like '%RunAtLoad%' and command like '%true%') and (command like '%LaunchAgents%' or command like '%LaunchDaemons%')",
  "Tags": [
    "attack.persistence",
    "attack.t1543.001",
    "attack.t1543.004"
  ],
  "Version": "0.1"
}