# csiem/linux/auditd/lnx_auditd_systemd_service_creation.json
{
  "Description": "Detects a creation of systemd services which could be used by adversaries to execute malicious code.\n\nAuthor: Pawel Mazur\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "1bac86ba-41aa-4f62-9d6b-405eac99b485",
  "InsertDate": "2025-01-21T09:27:13Z",
  "LastUpdateDate": "2025-01-21T09:27:13Z",
  "Level": "medium",
  "Name": "Systemd Service Creation",
  "Query": "sourcetype='linux-auditd' eql select * from _source_ where ((type = 'path' and nametype = 'create') and ((process_path like '/usr/lib/systemd/system/%' or process_path like '/etc/systemd/system/%') or process_path like '%/.config/systemd/user/%'))",
  "Tags": [
    "attack.persistence",
    "attack.t1543.002"
  ],
  "Version": "0.1"
}