# csiem/application/rpc_firewall/rpc_firewall_atsvc_lateral_movement.json
{
  "Description": "Detects remote RPC calls to create or execute a scheduled task via ATSvc\n\nAuthor: Sagie Dulce, Dekel Paz\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "0fcd1c79-4eeb-4746-aba9-1b458f7a79cb",
  "InsertDate": "2025-01-19T11:33:19Z",
  "LastUpdateDate": "2025-01-19T11:33:19Z",
  "Level": "high",
  "Name": "Remote Schedule Task Lateral Movement via ATSvc",
  "Query": "sourcetype like 'rpc_firewall-%' eql select * from _source_ where eventlog = 'rpcfw' and eventid = '3' and interfaceuuid = '1ff70682-0a51-30e8-076d-740be8cee98b' and (opnum = '0' or opnum = '1')",
  "Tags": [
    "attack.lateral-movement",
    "attack.t1053",
    "attack.t1053.002"
  ],
  "Version": "0.1"
}