# csiem/windows/builtin/taskscheduler/win_taskscheduler_execution_from_susp_locations.json
{
  "Description": "Detects the execution of Scheduled Tasks where the Program being run is located in a suspicious location or it's an unusale program to be run from a Scheduled Task\n\nAuthor: Nasreddine Bencherchali (Nextron Systems)\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "424273ea-7cf8-43a6-b712-375f925e481f",
  "InsertDate": "2025-01-25T21:20:40Z",
  "LastUpdateDate": "2025-01-25T21:20:40Z",
  "Level": "medium",
  "Name": "Scheduled Task Executed From A Suspicious Location",
  "Query": "sourcetype='windows-taskscheduler' eql select * from _source_ where eventid = '129' and (path like '%C:\\Windows\\Temp\\%' or path like '%\\AppData\\Local\\Temp\\%' or path like '%\\Desktop\\%' or path like '%\\Downloads\\%' or path like '%\\Users\\Public\\%' or path like '%C:\\Temp\\%')",
  "Tags": [
    "attack.persistence",
    "attack.t1053.005"
  ],
  "Version": "0.1"
}