# csiem/windows/builtin/taskscheduler/win_taskscheduler_lolbin_execution_via_task_scheduler.json
{
  "Description": "Detects the execution of Scheduled Tasks where the program being run is located in a suspicious location or where it is an unusual program to be run from a Scheduled Task\n\nAuthor: Nasreddine Bencherchali (Nextron Systems)\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "f0767f15-0fb3-44b9-851e-e8d9a6d0005d",
  "InsertDate": "2025-01-19T11:33:54Z",
  "LastUpdateDate": "2025-01-19T11:33:54Z",
  "Level": "medium",
  "Name": "Scheduled Task Executed Uncommon LOLBIN",
  "Query": "sourcetype='windows-taskscheduler' eql select * from _source_ where eventid = '129' and (path like '%\\calc.exe' or path like '%\\cscript.exe' or path like '%\\mshta.exe' or path like '%\\mspaint.exe' or path like '%\\notepad.exe' or path like '%\\regsvr32.exe' or path like '%\\wscript.exe')",
  "Tags": [
    "attack.persistence",
    "attack.t1053.005"
  ],
  "Version": "0.1"
}