# csiem/windows/process_creation/proc_creation_win_schtasks_appdata_local_system.json
{
  "Description": "Detects the creation of a schtask that executes a file from C:\\Users\\\u003cUSER\u003e\\AppData\\Local\n\nAuthor: pH-T (Nextron Systems), Nasreddine Bencherchali (Nextron Systems)\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "c5c00f49-b3f9-45a6-997e-cfdecc6e1967",
  "InsertDate": "2025-01-25T21:22:22Z",
  "LastUpdateDate": "2025-01-25T21:22:22Z",
  "Level": "high",
  "Name": "Suspicious Schtasks Execution AppData Folder",
  "Query": "sourcetype='windows-sysmon' eql select * from _source_ where ((process_path like '%\\schtasks.exe' and (command like '%/Create%' and command like '%/RU%' and command like '%/TR%' and command like '%C:\\Users\\%' and command like '%\\AppData\\Local\\%') and (command like '%NT AUT%' or command like '% SYSTEM %')) and  not ((parent_process_path like '%\\AppData\\Local\\Temp\\%' and parent_process_path like '%TeamViewer_.exe%') and process_path like '%\\schtasks.exe' and command like '%/TN TVInstallRestore%'))",
  "Tags": [
    "attack.execution",
    "attack.persistence",
    "attack.t1053.005",
    "attack.t1059.001"
  ],
  "Version": "0.1"
}