# csiem/windows/process_creation/proc_creation_win_schtasks_schedule_type.json
{
  "Description": "Detects scheduled task creations or modification on a suspicious schedule type\n\nAuthor: Nasreddine Bencherchali (Nextron Systems)\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "24c8392b-aa3c-46b7-a545-43f71657fe98",
  "InsertDate": "2025-01-25T21:22:23Z",
  "LastUpdateDate": "2025-01-25T21:22:23Z",
  "Level": "high",
  "Name": "Suspicious Schtasks Schedule Types",
  "Query": "sourcetype='windows-sysmon' eql select * from _source_ where (((process_path like '%\\schtasks.exe' and original_file_name = 'schtasks.exe') and (command like '% ONLOGON %' or command like '% ONSTART %' or command like '% ONCE %' or command like '% ONIDLE %')) and  not (command like '%NT AUT%' or command like '% SYSTEM%' or command like '%HIGHEST%'))",
  "Tags": [
    "attack.execution",
    "attack.t1053.005"
  ],
  "Version": "0.1"
}