# csiem/windows/process_creation/proc_creation_win_sc_service_path_modification.json
{
  "Description": "Detects service path modification via the \"sc\" binary to a suspicious command or path\n\nAuthor: Victor Sergeev, oscd.community, Nasreddine Bencherchali (Nextron Systems)\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "138d3531-8793-4f50-a2cd-f291b2863d78",
  "InsertDate": "2025-01-25T21:22:22Z",
  "LastUpdateDate": "2025-01-25T21:22:22Z",
  "Level": "high",
  "Name": "Suspicious Service Path Modification",
  "Query": "sourcetype='windows-sysmon' eql select * from _source_ where process_path like '%\\sc.exe' and (command like '%config%' and command like '%binPath%') and (command like '%powershell%' or command like '%cmd %' or command like '%mshta%' or command like '%wscript%' or command like '%cscript%' or command like '%rundll32%' or command like '%svchost%' or command like '%dllhost%' or command like '%cmd.exe /c%' or command like '%cmd.exe /k%' or command like '%cmd.exe /r%' or command like '%cmd /c%' or command like '%cmd /k%' or command like '%cmd /r%' or command like '%C:\\Users\\Public%' or command like '%\\Downloads\\%' or command like '%\\Desktop\\%' or command like '%\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\%' or command like '%C:\\Windows\\TEMP\\%' or command like '%\\AppData\\Local\\Temp%')",
  "Tags": [
    "attack.persistence",
    "attack.privilege-escalation",
    "attack.t1543.003"
  ],
  "Version": "0.1"
}