# csiem/windows/process_creation/proc_creation_win_schtasks_guid_task_name.json
{
  "Description": "Detects creation of a scheduled task with a GUID like name\n\nAuthor: Nasreddine Bencherchali (Nextron Systems)\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "ff2fff64-4cd6-4a2b-ba7d-e28a30bbe66b",
  "InsertDate": "2025-01-19T11:35:26Z",
  "LastUpdateDate": "2025-01-19T11:35:26Z",
  "Level": "medium",
  "Name": "Suspicious Scheduled Task Name As GUID",
  "Query": "sourcetype='windows-sysmon' eql select * from _source_ where ((process_path like '%\\schtasks.exe' and command like '%/create %') and (command like '%/tn \"{%' or command like '%/tn '{%' or command like '%/tn {%') and (command like '%}\"%' or command like '%}'%' or command like '%} %'))",
  "Tags": [
    "attack.execution",
    "attack.t1053.005"
  ],
  "Version": "0.1"
}