# csiem/windows/process_creation/proc_creation_win_schtasks_schedule_via_masqueraded_xml_file.json
{
  "Description": "Detects the creation of a scheduled task using the \"-XML\" flag with a file without the '.xml' extension. This behavior could be indicative of potential defense evasion attempt during persistence\n\nAuthor: Swachchhanda Shrawan Poudel, Elastic (idea)\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "dd2a821e-3b07-4d3b-a9ac-929fe4c6ca0c",
  "InsertDate": "2025-01-25T21:22:23Z",
  "LastUpdateDate": "2025-01-25T21:22:23Z",
  "Level": "medium",
  "Name": "Suspicious Scheduled Task Creation via Masqueraded XML File",
  "Query": "sourcetype='windows-sysmon' eql select * from _source_ where (((process_path like '%\\schtasks.exe' and original_file_name = 'schtasks.exe') and (command like '%/create%' or command like '%-create%') and (command like '%/xml%' or command like '%-xml%')) and  not (command like '%.xml%' or (integrity_level = 'System' or integrity_level = 'S-1-16-16384') or (parent_process_path like '%\\rundll32.exe' and (parent_command_line like '%:\\WINDOWS\\Installer\\MSI%' and parent_command_line like '%.tmp,zzzzInvokeManagedCustomActionOutOfProc%'))) and  not (parent_process_path like '%:\\ProgramData\\OEM\\UpgradeTool\\CareCenter_*\\BUnzip\\Setup_msi.exe' or parent_process_path like '%:\\Program Files\\Axis Communications\\AXIS Camera Station\\SetupActions.exe' or parent_process_path like '%:\\Program Files\\Axis Communications\\AXIS Device Manager\\AdmSetupActions.exe' or parent_process_path like '%:\\Program Files (x86)\\Zemana\\AntiMalware\\AntiMalware.exe' or parent_process_path like '%:\\Program Files\\Dell\\SupportAssist\\pcdrcui.exe'))",
  "Tags": [
    "attack.defense-evasion",
    "attack.persistence",
    "attack.t1036.005",
    "attack.t1053.005"
  ],
  "Version": "0.1"
}