# csiem/windows/registry/registry_event/registry_event_runkey_winekey.json
{
  "Description": "Detects potential malicious modification of run keys by winekey or team9 backdoor\n\nAuthor: omkar72\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "b98968aa-dbc0-4a9c-ac35-108363cbf8d5",
  "InsertDate": "2025-01-25T21:22:41Z",
  "LastUpdateDate": "2025-01-25T21:22:41Z",
  "Level": "high",
  "Name": "WINEKEY Registry Modification",
  "Query": "sourcetype='windows-sysmon' eql select * from _source_ where target_object like '%Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Backup Mgr'",
  "Tags": [
    "attack.persistence",
    "attack.t1547"
  ],
  "Version": "0.1"
}