# splunk/windows/file/file_event/file_event_win_ripzip_attack.json
{
  "Description": "Detects a phishing attack which expands a ZIP file containing a malicious shortcut.\nIf the victim expands the ZIP file via the explorer process, then the explorer process expands the malicious ZIP file and drops a malicious shortcut redirected to a backdoor into the Startup folder.\nAdditionally, the file name of the malicious shortcut in Startup folder contains {0AFACED1-E828-11D1-9187-B532F1E9575D} meaning the folder shortcut operation.\n\n\nAuthor: Greg (rule)\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "a6976974-ea6f-4e97-818e-ea08625c52cb",
  "InsertDate": "2025-02-03T20:32:19Z",
  "LastUpdateDate": "2025-02-03T20:32:19Z",
  "Level": "high",
  "Name": "Potential RipZip Attack on Startup Folder",
  "Query": "sourcetype=\"windows-sysmon\" (filepath=\"*\\\\Microsoft\\\\Windows\\\\Start Menu\\\\Programs\\\\Startup*\" filepath=\"*.lnk.{0AFACED1-E828-11D1-9187-B532F1E9575D}*\") process_path=\"*\\\\explorer.exe\"",
  "Tags": [
    "attack.persistence",
    "attack.t1547"
  ],
  "Version": "0.1"
}