# csiem/windows/builtin/security/win_security_susp_scheduled_task_creation.json
{
  "Description": "Detects suspicious scheduled task creation events. Based on attributes such as paths, commands line flags, etc.\n\nAuthor: Nasreddine Bencherchali (Nextron Systems)\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "3a734d25-df5c-4b99-8034-af1ddb5883a4",
  "InsertDate": "2025-01-19T11:33:50Z",
  "LastUpdateDate": "2025-01-19T11:33:50Z",
  "Level": "high",
  "Name": "Suspicious Scheduled Task Creation",
  "Query": "sourcetype='windows-security' eql select * from _source_ where (eventid = '4698' and (taskcontent like '%\\appdata\\local\\temp\\%' or taskcontent like '%\\appdata\\roaming\\%' or taskcontent like '%\\users\\public\\%' or taskcontent like '%\\windows\\temp\\%' or taskcontent like '%c:\\temp\\%' or taskcontent like '%\\desktop\\%' or taskcontent like '%\\downloads\\%' or taskcontent like '%\\temporary internet%' or taskcontent like '%c:\\programdata\\%' or taskcontent like '%c:\\perflogs\\%') and (taskcontent like '%regsvr32%' or taskcontent like '%rundll32%' or taskcontent like '%cmd.exe\u003c/command\u003e%' or taskcontent like '%cmd\u003c/command\u003e%' or taskcontent like '%\u003carguments\u003e/c %' or taskcontent like '%\u003carguments\u003e/k %' or taskcontent like '%\u003carguments\u003e/r %' or taskcontent like '%powershell%' or taskcontent like '%pwsh%' or taskcontent like '%mshta%' or taskcontent like '%wscript%' or taskcontent like '%cscript%' or taskcontent like '%certutil%' or taskcontent like '%bitsadmin%' or taskcontent like '%bash.exe%' or taskcontent like '%bash %' or taskcontent like '%scrcons%' or taskcontent like '%wmic %' or taskcontent like '%wmic.exe%' or taskcontent like '%forfiles%' or taskcontent like '%scriptrunner%' or taskcontent like '%hh.exe%'))",
  "Tags": [
    "attack.execution",
    "attack.privilege-escalation",
    "attack.persistence",
    "attack.t1053.005"
  ],
  "Version": "0.1"
}