# csiem/windows/builtin/security/win_security_susp_scheduled_task_delete_or_disable.json
{
  "Description": "Detects when adversaries stop services or processes by deleting or disabling their respective scheduled tasks in order to conduct data destructive activities\n\nAuthor: Nasreddine Bencherchali (Nextron Systems)\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "7595ba94-cf3b-4471-aa03-4f6baa9e5fad",
  "InsertDate": "2025-01-25T21:20:27Z",
  "LastUpdateDate": "2025-01-25T21:20:27Z",
  "Level": "high",
  "Name": "Important Scheduled Task Deleted/Disabled",
  "Query": "sourcetype='windows-security' eql select * from _source_ where (((eventid = '4699' or eventid = '4701') and (taskname like '%\\Windows\\SystemRestore\\SR%' or taskname like '%\\Windows\\Windows Defender\\%' or taskname like '%\\Windows\\BitLocker%' or taskname like '%\\Windows\\WindowsBackup\\%' or taskname like '%\\Windows\\WindowsUpdate\\%' or taskname like '%\\Windows\\UpdateOrchestrator\\Schedule%' or taskname like '%\\Windows\\ExploitGuard%')) and  not (eventid = '4699' and subjectusername like '%$' and taskname like '%\\Windows\\Windows Defender\\%'))",
  "Tags": [
    "attack.execution",
    "attack.privilege-escalation",
    "attack.persistence",
    "attack.t1053.005"
  ],
  "Version": "0.1"
}