# csiem/windows/builtin/security/win_security_susp_scheduled_task_update.json
{
  "Description": "Detects update to a scheduled task event that contain suspicious keywords.\n\nAuthor: Nasreddine Bencherchali (Nextron Systems)\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "614cf376-6651-47c4-9dcc-6b9527f749f4",
  "InsertDate": "2025-01-25T21:20:27Z",
  "LastUpdateDate": "2025-01-25T21:20:27Z",
  "Level": "high",
  "Name": "Suspicious Scheduled Task Update",
  "Query": "sourcetype='windows-security' eql select * from _source_ where ((taskcontentnew like '%\\AppData\\Local\\Temp\\%' or taskcontentnew like '%\\AppData\\Roaming\\%' or taskcontentnew like '%\\Users\\Public\\%' or taskcontentnew like '%\\WINDOWS\\Temp\\%' or taskcontentnew like '%C:\\Temp\\%' or taskcontentnew like '%\\Desktop\\%' or taskcontentnew like '%\\Downloads\\%' or taskcontentnew like '%\\Temporary Internet%' or taskcontentnew like '%C:\\ProgramData\\%' or taskcontentnew like '%C:\\Perflogs\\%') and (taskcontentnew like '%regsvr32%' or taskcontentnew like '%rundll32%' or taskcontentnew like '%cmd.exe\u003c/Command\u003e%' or taskcontentnew like '%cmd\u003c/Command\u003e%' or taskcontentnew like '%\u003cArguments\u003e/c %' or taskcontentnew like '%\u003cArguments\u003e/k %' or taskcontentnew like '%\u003cArguments\u003e/r %' or taskcontentnew like '%powershell%' or taskcontentnew like '%pwsh%' or taskcontentnew like '%mshta%' or taskcontentnew like '%wscript%' or taskcontentnew like '%cscript%' or taskcontentnew like '%certutil%' or taskcontentnew like '%bitsadmin%' or taskcontentnew like '%bash.exe%' or taskcontentnew like '%bash %' or taskcontentnew like '%scrcons%' or taskcontentnew like '%wmic %' or taskcontentnew like '%wmic.exe%' or taskcontentnew like '%forfiles%' or taskcontentnew like '%scriptrunner%' or taskcontentnew like '%hh.exe%') and eventid = '4702')",
  "Tags": [
    "attack.execution",
    "attack.privilege-escalation",
    "attack.persistence",
    "attack.t1053.005"
  ],
  "Version": "0.1"
}