# csiem/windows/builtin/security/win_security_susp_scheduled_task_creation.json
{
  "Description": "Detects suspicious scheduled task creation events. Based on attributes such as paths, commands line flags, etc.\n\nAuthor: Nasreddine Bencherchali (Nextron Systems)\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "3a734d25-df5c-4b99-8034-af1ddb5883a4",
  "InsertDate": "2025-01-25T21:20:27Z",
  "LastUpdateDate": "2025-01-25T21:20:27Z",
  "Level": "high",
  "Name": "Suspicious Scheduled Task Creation",
  "Query": "sourcetype='windows-security' eql select * from _source_ where (eventid = '4698' and (taskcontent like '%\\AppData\\Local\\Temp\\%' or taskcontent like '%\\AppData\\Roaming\\%' or taskcontent like '%\\Users\\Public\\%' or taskcontent like '%\\WINDOWS\\Temp\\%' or taskcontent like '%C:\\Temp\\%' or taskcontent like '%\\Desktop\\%' or taskcontent like '%\\Downloads\\%' or taskcontent like '%\\Temporary Internet%' or taskcontent like '%C:\\ProgramData\\%' or taskcontent like '%C:\\Perflogs\\%') and (taskcontent like '%regsvr32%' or taskcontent like '%rundll32%' or taskcontent like '%cmd.exe\u003c/Command\u003e%' or taskcontent like '%cmd\u003c/Command\u003e%' or taskcontent like '%\u003cArguments\u003e/c %' or taskcontent like '%\u003cArguments\u003e/k %' or taskcontent like '%\u003cArguments\u003e/r %' or taskcontent like '%powershell%' or taskcontent like '%pwsh%' or taskcontent like '%mshta%' or taskcontent like '%wscript%' or taskcontent like '%cscript%' or taskcontent like '%certutil%' or taskcontent like '%bitsadmin%' or taskcontent like '%bash.exe%' or taskcontent like '%bash %' or taskcontent like '%scrcons%' or taskcontent like '%wmic %' or taskcontent like '%wmic.exe%' or taskcontent like '%forfiles%' or taskcontent like '%scriptrunner%' or taskcontent like '%hh.exe%'))",
  "Tags": [
    "attack.execution",
    "attack.privilege-escalation",
    "attack.persistence",
    "attack.t1053.005"
  ],
  "Version": "0.1"
}