# csiem/windows/registry/registry_set/registry_set_asep_reg_keys_modification_wow6432node.json
{
  "Description": "Detects modification of autostart extensibility point (ASEP) in registry.\n\nAuthor: Victor Sergeev, Daniil Yugoslavskiy, Gleb Sukhodolskiy, Timur Zinniatullin, oscd.community, Tim Shelton, frack113 (split)\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "b29aed60-ebd1-442b-9cb5-16a1d0324adb",
  "InsertDate": "2025-01-25T21:22:43Z",
  "LastUpdateDate": "2025-01-25T21:22:43Z",
  "Level": "medium",
  "Name": "Wow6432Node CurrentVersion Autorun Keys Modification",
  "Query": "sourcetype like 'windows-%' eql select * from _source_ where ((target_object like '%\\SOFTWARE\\Wow6432Node\\Microsoft\\Windows\\CurrentVersion%' and (target_object like '%\\ShellServiceObjectDelayLoad%' or target_object like '%\\Run\\%' or target_object like '%\\RunOnce\\%' or target_object like '%\\RunOnceEx\\%' or target_object like '%\\RunServices\\%' or target_object like '%\\RunServicesOnce\\%' or target_object like '%\\Explorer\\ShellServiceObjects%' or target_object like '%\\Explorer\\ShellIconOverlayIdentifiers%' or target_object like '%\\Explorer\\ShellExecuteHooks%' or target_object like '%\\Explorer\\SharedTaskScheduler%' or target_object like '%\\Explorer\\Browser Helper Objects%')) and  not (details like '\"C:\\ProgramData\\Package Cache\\{d21a4f20-968a-4b0c-bf04-a38da5f06e41}\\windowsdesktop-runtime-%' or (process_path like '%\\VC_redist.x64.exe' and details like '%}\\VC_redist.x64.exe\" /burn.runonce') or (process_path = 'C:\\Program Files\\Common Files\\Microsoft Shared\\ClickToRun\\OfficeClickToRun.exe' and target_object like '%\\Office\\ClickToRun\\REGISTRY\\MACHINE\\Software\\Wow6432Node\\%') or target_object like '%\\Explorer\\Browser Helper Objects\\{92EF2EAD-A7CE-4424-B0DB-499CF856608E}\\NoExplorer' or (details like '%-A251-47B7-93E1-CDD82E34AF8B}' and details = 'grpconv -o' and (details like '%C:\\Program Files%' and details like '%\\Dropbox\\Client\\Dropbox.exe%' and details like '% /systemstartup%')) or (process_path like '%\\windowsdesktop-runtime-%' and (target_object like '%\\WOW6432Node\\Microsoft\\Windows\\CurrentVersion\\RunOnce\\{e2d1ae32-dd1d-4ad7-a298-10e42e7840fc}' or target_object like '%\\WOW6432Node\\Microsoft\\Windows\\CurrentVersion\\RunOnce\\{7037b699-7382-448c-89a7-4765961d2537}') and details like '\"C:\\ProgramData\\Package Cache\\%' and details like '%.exe\" /burn.runonce') or ((process_path like 'C:\\Program Files\\Common Files\\Microsoft Shared\\ClickToRun\\%' or process_path like 'C:\\Program Files\\Common Files\\Microsoft Shared\\ClickToRun\\Updates\\%') and process_path like '%\\OfficeClickToRun.exe') or (process_path = 'C:\\WINDOWS\\system32\\msiexec.exe' and target_object like '%\\SOFTWARE\\Wow6432Node\\Microsoft\\Windows\\CurrentVersion\\Run\\%') or details = '(Empty)' or (process_path like '%C:\\Program Files (x86)\\Microsoft\\EdgeUpdate\\Install\\{%' and process_path like '%\\setup.exe%') or (process_path like 'C:\\Windows\\Installer\\MSI%' and target_object like '%\\Explorer\\Browser Helper Objects%') or ((process_path = 'C:\\Program Files\\Microsoft Office\\root\\integration\\integrator.exe' or process_path = 'C:\\Program Files (x86)\\Microsoft Office\\root\\integration\\integrator.exe') and target_object like '%\\Explorer\\Browser Helper Objects\\{31D09BA0-12F5-4CCE-BE8A-2923E76605DA}\\%') or ((process_path like 'C:\\ProgramData\\Package Cache%' or process_path like 'C:\\Windows\\Temp\\%') and (process_path like '%\\winsdksetup.exe%' or process_path like '%\\windowsdesktop-runtime-%' or process_path like '%\\AspNetCoreSharedFrameworkBundle-%') and details like '% /burn.runonce')))",
  "Tags": [
    "attack.persistence",
    "attack.t1547.001"
  ],
  "Version": "0.1"
}