# csiem/windows/registry/registry_set/registry_set_powershell_in_run_keys.json
{
  "Description": "Detects potential PowerShell commands or code within registry run keys\n\nAuthor: frack113, Florian Roth (Nextron Systems)\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "8d85cf08-bf97-4260-ba49-986a2a65129c",
  "InsertDate": "2025-01-25T21:22:49Z",
  "LastUpdateDate": "2025-01-25T21:22:49Z",
  "Level": "medium",
  "Name": "Suspicious Powershell In Registry Run Keys",
  "Query": "sourcetype like 'windows-%' eql select * from _source_ where target_object like '%\\Software\\Microsoft\\Windows\\CurrentVersion\\Run%' and (details like '%powershell%' or details like '%pwsh %' or details like '%FromBase64String%' or details like '%.DownloadFile(%' or details like '%.DownloadString(%' or details like '% -w hidden %' or details like '% -w 1 %' or details like '%-windowstyle hidden%' or details like '%-window hidden%' or details like '% -nop %' or details like '% -encodedcommand %' or details like '%-ExecutionPolicy Bypass%' or details like '%Invoke-Expression%' or details like '%IEX (%' or details like '%Invoke-Command%' or details like '%ICM -%' or details like '%Invoke-WebRequest%' or details like '%IWR %' or details like '% -noni %' or details like '% -noninteractive %')",
  "Tags": [
    "attack.persistence",
    "attack.t1547.001"
  ],
  "Version": "0.1"
}