# csiem/windows/registry/registry_set/registry_set_asep_reg_keys_modification_currentversion.json
{
  "Description": "Detects modification of autostart extensibility point (ASEP) in registry.\n\nAuthor: Victor Sergeev, Daniil Yugoslavskiy, Gleb Sukhodolskiy, Timur Zinniatullin, oscd.community, Tim Shelton, frack113 (split)\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "20f0ee37-5942-4e45-b7d5-c5b5db9df5cd",
  "InsertDate": "2025-01-25T21:22:42Z",
  "LastUpdateDate": "2025-01-25T21:22:42Z",
  "Level": "medium",
  "Name": "CurrentVersion Autorun Keys Modification",
  "Query": "sourcetype like 'windows-%' eql select * from _source_ where (((target_object like '%\\ShellServiceObjectDelayLoad%' or target_object like '%\\Run\\%' or target_object like '%\\RunOnce\\%' or target_object like '%\\RunOnceEx\\%' or target_object like '%\\RunServices\\%' or target_object like '%\\RunServicesOnce\\%' or target_object like '%\\Policies\\System\\Shell%' or target_object like '%\\Policies\\Explorer\\Run%' or target_object like '%\\Group Policy\\Scripts\\Startup%' or target_object like '%\\Group Policy\\Scripts\\Shutdown%' or target_object like '%\\Group Policy\\Scripts\\Logon%' or target_object like '%\\Group Policy\\Scripts\\Logoff%' or target_object like '%\\Explorer\\ShellServiceObjects%' or target_object like '%\\Explorer\\ShellIconOverlayIdentifiers%' or target_object like '%\\Explorer\\ShellExecuteHooks%' or target_object like '%\\Explorer\\SharedTaskScheduler%' or target_object like '%\\Explorer\\Browser Helper Objects%' or target_object like '%\\Authentication\\PLAP Providers%' or target_object like '%\\Authentication\\Credential Providers%' or target_object like '%\\Authentication\\Credential Provider Filters%') and target_object like '%\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion%') and  not ((process_path like 'C:\\Program Files\\AVG\\Antivirus\\Setup\\%' and (details = '\"C:\\Program Files\\AVG\\Antivirus\\AvLaunch.exe\" /gui' or details = '\"C:\\Program Files (x86)\\AVG\\Antivirus\\AvLaunch.exe\" /gui' or details = '{472083B0-C522-11CF-8763-00608CC02F24}')) or (target_object like '%\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\Opera Browser Assistant' and details = 'C:\\Program Files\\Opera\\assistant\\browser_assistant.exe') or (target_object like '%\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\Greenshot' and details = 'C:\\Program Files\\Greenshot\\Greenshot.exe') or (target_object like '%GoogleDrive%' and (details = '{CFE8B367-77A7-41D7-9C90-75D16D7DC6B6}' or details = '{A8E52322-8734-481D-A7E2-27B309EF8D56}' or details = '{C973DA94-CBDF-4E77-81D1-E5B794FBD146}' or details = '{51EF1569-67EE-4AD6-9646-E726C3FFC8A2}')) or (process_path like '%\\Microsoft\\Teams\\current\\Teams.exe' and details like '%\\Microsoft\\Teams\\Update.exe --processStart %') or (process_path = 'C:\\Windows\\system32\\userinit.exe' and details = 'ctfmon.exe /n') or (process_path like 'C:\\Program Files (x86)\\Microsoft\\EdgeUpdate\\Install\\%' or process_path like 'C:\\Program Files (x86)\\Microsoft\\EdgeWebView\\%' or process_path like 'C:\\Program Files (x86)\\Microsoft\\Edge\\Application\\msedge.exe%') or (process_path = 'C:\\Windows\\system32\\LogonUI.exe' and (target_object like '%\\Authentication\\Credential Providers\\{D6886603-9D2F-4EB2-B667-1971041FA96B}\\%' or target_object like '%\\Authentication\\Credential Providers\\{BEC09223-B018-416D-A0AC-523971B639F5}\\%' or target_object like '%\\Authentication\\Credential Providers\\{8AF662BF-65A0-4D0A-A540-A338A999D36F}\\%' or target_object like '%\\Authentication\\Credential Providers\\{27FBDB57-B613-4AF2-9D7E-4FA7A66C21AD}\\%')) or ((details like 'C:\\Windows\\system32\\cmd.exe /q /c rmdir /s /q \"C:\\Users\\%' or details like 'C:\\Windows\\system32\\cmd.exe /q /c del /q \"C:\\Users\\%') and details like '%\\AppData\\Local\\Microsoft\\OneDrive\\%') or ((process_path like 'C:\\Program Files\\Common Files\\Microsoft Shared\\ClickToRun\\%' or process_path like 'C:\\Program Files\\Common Files\\Microsoft Shared\\ClickToRun\\Updates\\%') and process_path like '%\\OfficeClickToRun.exe') or ((process_path like '%\\aurora-agent-64.exe' or process_path like '%\\aurora-agent.exe') and target_object like '%\\Microsoft\\Windows\\CurrentVersion\\Run\\aurora-dashboard' and details = 'C:\\Program Files\\Aurora-Agent\\tools\\aurora-dashboard.exe') or (details = '(Empty)' and target_object like '%\\NgcFirst\\ConsecutiveSwitchCount' and (process_path like '%\\AppData\\Local\\Microsoft\\OneDrive\\Update\\OneDriveSetup.exe' or process_path like '%\\AppData\\Roaming\\Spotify\\Spotify.exe' or process_path like '%\\AppData\\Local\\WebEx\\WebexHost.exe') and (process_path = 'C:\\WINDOWS\\system32\\devicecensus.exe' or process_path = 'C:\\Windows\\system32\\winsat.exe' or process_path = 'C:\\Program Files\\Microsoft OneDrive\\StandaloneUpdater\\OneDriveSetup.exe' or process_path = 'C:\\Program Files\\Microsoft OneDrive\\Update\\OneDriveSetup.exe' or process_path = 'C:\\Program Files (x86)\\Microsoft OneDrive\\Update\\OneDriveSetup.exe' or process_path = 'C:\\Program Files\\KeePass Password Safe 2\\ShInstUtil.exe' or process_path = 'C:\\Program Files\\Everything\\Everything.exe' or process_path = 'C:\\Program Files (x86)\\Microsoft Office\\root\\integration\\integrator.exe')) or (process_path = 'C:\\Windows\\system32\\regsvr32.exe' and target_object like '%DropboxExt%' and details like '%A251-47B7-93E1-CDD82E34AF8B}') or (target_object like '%\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\iTunesHelper' and details = '\"C:\\Program Files\\iTunes\\iTunesHelper.exe\"') or (target_object like '%\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\RunOnce\\zoommsirepair' and details = '\"C:\\Program Files\\Zoom\\bin\\installer.exe\" /repair') or (target_object like '%\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\GoogleDriveFS' and details like 'C:\\Program Files\\Google\\Drive File Stream\\%' and details like '%\\GoogleDriveFS.exe%') or (target_object like '%\\Microsoft\\Windows\\CurrentVersion\\RunOnce\\{%' and (details like '%\\AppData\\Local\\Package Cache\\{%' and details like '%}\\python-%') and details like '%.exe\" /burn.runonce') or process_path = 'C:\\Program Files\\Windows Defender\\MsMpEng.exe' or (target_object like '%\\Microsoft\\Windows\\CurrentVersion\\Run\\Everything' and details like '%\\Everything\\Everything.exe\" -startup')))",
  "Tags": [
    "attack.persistence",
    "attack.t1547.001"
  ],
  "Version": "0.1"
}