# csiem/windows/registry/registry_set/registry_set_vbs_payload_stored.json
{
  "Description": "Detects VBScript content stored into registry keys as seen being used by UNC2452 group\n\nAuthor: Florian Roth (Nextron Systems)\nSigma Repository: [GitHub](https://github.com/SigmaHQ/sigma)",
  "ID": "46490193-1b22-4c29-bdd6-5bf63907216f",
  "InsertDate": "2025-01-25T21:22:51Z",
  "LastUpdateDate": "2025-01-25T21:22:51Z",
  "Level": "high",
  "Name": "VBScript Payload Stored in Registry",
  "Query": "sourcetype like 'windows-%' eql select * from _source_ where ((target_object like '%Software\\Microsoft\\Windows\\CurrentVersion%' and (details like '%vbscript:%' or details like '%jscript:%' or details like '%mshtml,%' or details like '%RunHTMLApplication%' or details like '%Execute(%' or details like '%CreateObject%' or details like '%window.close%')) and  not (target_object like '%Software\\Microsoft\\Windows\\CurrentVersion\\Run%' or (process_path like '%\\msiexec.exe' and target_object like '%\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Installer\\UserData\\%' and (details like '%\\Microsoft.NET\\Primary Interop Assemblies\\Microsoft.mshtml.dll%' or details like '%\u003c\\Microsoft.mshtml,fileVersion=%' or details like '%_mshtml_dll_%' or details like '%\u003c\\Microsoft.mshtml,culture=%'))))",
  "Tags": [
    "attack.persistence",
    "attack.t1547.001"
  ],
  "Version": "0.1"
}